# streamlit/deploy/weather-streamlit.service
# /etc/systemd/system/weather-streamlit.service
#
# Instalação:
#   sudo cp deploy/weather-streamlit.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable weather-streamlit
#   sudo systemctl start weather-streamlit

[Unit]
Description=Weather Analytics Streamlit App
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Analytics-Engineer/Weather-Analytics/streamlit
EnvironmentFile=/home/ubuntu/Analytics-Engineer/Weather-Analytics/streamlit/.env

# Ajuste o caminho do venv se necessário (ver seção de deploy no README)
ExecStart=/home/ubuntu/venv/weather/bin/streamlit run app.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --browser.gatherUsageStats=false

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-streamlit

[Install]
WantedBy=multi-user.target
